# roles/elk_config/templates/logstash.service.j2
[Unit]
Description=Logstash service
After=network.target

[Service]
Type=simple
ExecStart=/data/{{ logstash }}/start.sh

[Install]
WantedBy=multi-user.target
