# roles/setup_ELK/templates/kibana.service.j2
[Unit]
Description=Kibana service
After=network.target

[Service]
Type=simple
ExecStart=/data/kibana-6.4.0-linux-x86_64/bin/kibana

[Install]
WantedBy=multi-user.target
